# dvk512/lcd_monitor.cron
* * * * * root /usr/local/bin/lcd1602_monitor
