# dkv512/led_clock.cron
* * * * * root /usr/local/bin/led_clock.sh
